# deploy/beforeu-api.service
[Unit]
Description=BeforeUgone Phase 1 API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=beforeu-api
Group=beforeu-api
EnvironmentFile=/etc/beforeu-api.env
WorkingDirectory=/var/lib/beforeu-api
ExecStart=/usr/local/bin/beforeu-api
Restart=on-failure
RestartSec=3s
TimeoutStopSec=20s
KillSignal=SIGTERM
UMask=0077

StateDirectory=beforeu-api
StateDirectoryMode=0750
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=
MemoryMax=192M
TasksMax=64
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
