# deploy/beforeu-api-backup.timer
[Unit]
Description=Daily BeforeUgone API backup

[Timer]
OnCalendar=*-*-* 03:15:00
RandomizedDelaySec=15m
Persistent=true
Unit=beforeu-api-backup.service

[Install]
WantedBy=timers.target
